# agent/docker/serf/supervisor_conf/supervisord-restart.conf
[program:restart-proc]
command=/restart-proc.py
numprocs=1
autostart=true
autorestart=true
