# agent/docker/serf/supervisor_conf/supervisord-memory.conf
[program:memory-simulator-proc]
command=/memory-simulator-proc.sh
numprocs=1
autostart=true
autorestart=true
